# systemd/.config/systemd/user/anki-sync-server.service
[Unit]
Description=Anki Sync Server

[Service]
ExecStart=%h/.dotfiles/nopath-scripts/run_anki_sync.sh
